# roles/activemq/templates/amq_broker.d.service.j2
[Service]
ExecStopPost=usermod --home {{ activemq.user_home }} {{ activemq_service_user }}
